# scripts/crontab
* * * * * www php /var/www/backend/artisan schedule:run > /dev/null 2>&1